# com.secondbrain.ui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.secondbrain.ui</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/brentrossin/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>-m</string>
        <string>secondbrain.ui</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/brentrossin/SecondBrain</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/brentrossin/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/brentrossin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/brentrossin/SecondBrain/data/ui.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/brentrossin/SecondBrain/data/ui.log</string>
</dict>
</plist>
